# cmake/templates/app.service.in
[Unit]
Description=@APP_DISPLAY_NAME@
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/@PROJECT_NAME@
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
